# curator4-2144483268.conf
[program:curator4-2144483268]
command=python3 /home/curator4/curator4-2144483268/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-2144483268/test.err.log
stdout_logfile=/home/curator4/curator4-2144483268/test.out.log

